# v2-reference/uDOS-host/config/systemd/udos-hostd.service
[Unit]
Description=uDOS host bootstrap supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/udos/udos.env
ExecStart=/usr/bin/env bash -lc '"${UDOS_INSTALL_ROOT:-/opt/udos}/scripts/udos-hostd.sh"'
Restart=on-failure

[Install]
WantedBy=multi-user.target
