# v2-reference/uDOS-host/config/systemd/udos-gitd.service
[Unit]
Description=uDOS git repo-store and GitHub host worker
After=udos-commandd.service udos-syncd.service
Requires=udos-commandd.service

[Service]
Type=simple
EnvironmentFile=/etc/udos/udos.env
ExecStart=/usr/bin/env bash -lc '"${UDOS_INSTALL_ROOT:-/opt/udos}/scripts/udos-gitd.sh"'
Restart=on-failure

[Install]
WantedBy=multi-user.target
