# v2-reference/uDOS-host/config/systemd/udos-networkd.service
[Unit]
Description=uDOS network host, beacon, and portal service
After=udos-commandd.service
Requires=udos-commandd.service

[Service]
Type=simple
EnvironmentFile=/etc/udos/udos.env
EnvironmentFile=-/etc/udos/udos-networkd.env
ExecStart=/usr/bin/env bash -lc '"${UDOS_INSTALL_ROOT:-/opt/udos}/scripts/udos-networkd.sh"'
Restart=on-failure

[Install]
WantedBy=multi-user.target
