# v2-reference/uDOS-host/config/systemd/udos-web.service
[Unit]
Description=uDOS browser command centre
After=udos-commandd.service udos-vaultd.service udos-syncd.service udos-scheduled.service udos-networkd.service udos-budgetd.service
Requires=udos-commandd.service

[Service]
Type=simple
EnvironmentFile=/etc/udos/udos.env
EnvironmentFile=-/etc/udos/udos-web.env
ExecStart=/usr/bin/env bash -lc '"${UDOS_INSTALL_ROOT:-/opt/udos}/scripts/udos-web.sh"'
Restart=on-failure

[Install]
WantedBy=multi-user.target
