# v2-reference/uDOS-host/config/systemd/udos-commandd.service
[Unit]
Description=uDOS command router and operation registry
After=udos-hostd.service
Requires=udos-hostd.service

[Service]
Type=simple
EnvironmentFile=/etc/udos/udos.env
ExecStart=/usr/bin/env bash -lc '"${UDOS_INSTALL_ROOT:-/opt/udos}/scripts/udos-commandd.sh"'
Restart=on-failure

[Install]
WantedBy=multi-user.target
