# v2-reference/uDOS-host/config/systemd/udos-syncd.service
[Unit]
Description=uDOS sync queue and replication worker
After=udos-commandd.service udos-vaultd.service
Requires=udos-commandd.service

[Service]
Type=simple
EnvironmentFile=/etc/udos/udos.env
ExecStart=/usr/bin/env bash -lc '"${UDOS_INSTALL_ROOT:-/opt/udos}/scripts/udos-syncd.sh"'
Restart=on-failure

[Install]
WantedBy=multi-user.target
