# v2-reference/uDOS-host/config/systemd/udos-wizard-adapter.service
[Unit]
Description=uDOS Wizard adapter services
After=udos-networkd.service udos-budgetd.service
Wants=udos-networkd.service udos-budgetd.service

[Service]
Type=simple
EnvironmentFile=/etc/udos/udos.env
EnvironmentFile=-/etc/udos/udos-wizard-adapter.env
ExecStart=/usr/bin/env bash -lc '"${UDOS_INSTALL_ROOT:-/opt/udos}/scripts/udos-wizard-adapter.sh"'
Restart=on-failure

[Install]
WantedBy=multi-user.target
